# ai-voice-assistant.service
[Unit]
Description=AI Voice Assistant
After=graphical-session.target pipewire.service pulseaudio.service
Wants=pipewire.service

[Service]
Type=simple
User=%%USER%%
Environment=DISPLAY=:0
Environment=XDG_CURRENT_DESKTOP=%%DESKTOP%%
Environment=WAYLAND_DISPLAY=%%WAYLAND%%
Environment=XDG_RUNTIME_DIR=/run/user/%%UID%%
Environment=PATH=%%HOME%%/.opencode/bin:/usr/local/bin:/usr/bin:/bin:%%HOME%%/.local/bin
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%%UID%%/bus
ExecStart=%%HOME%%/.local/share/ai-voice-assistant/venv/bin/python %%HOME%%/.local/share/ai-voice-assistant/main.py
WorkingDirectory=%%HOME%%/.local/share/ai-voice-assistant
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
